# build/master.project/ignition/dropins/git-daemon.service.d/50-path.conf
[Unit]
ConditionPathExists=/srv/.git

[Service]
Environment="Path=/srv"
