# build/master.project/ignition/units/git-daemon.service
[Unit]
Description=git-daemon
Documentation=man:git-daemon(1)

[Service]
Restart=always
RestartSec=10s

ExecStart=/usr/libexec/git-core/git-daemon \
    --syslog --verbose --informative-errors \
    --export-all --enable=upload-pack --enable=upload-archive --enable=receive-pack \
    ${Path}

[Install]
WantedBy=multi-user.target
