# build/master.project/ignition/units/setup-network-environment.service
[Unit]
Description=Setup Network Environment
Documentation=https://github.com/kelseyhightower/setup-network-environment
Requires=network-online.target
After=network-online.target

[Service]
Type=oneshot
#ExecStartPre=-/bin/mkdir -p ${LIBEXEC}"
#ExecStartPre=/bin/wget -N -P ${LIBEXEC} ${URL}
#ExecStartPre=/bin/chmod +x ${LIBEXEC}/setup-network-environment

ExecStart=/bin/sh -c \
    "${LIBEXEC}/setup-network-environment"
RemainAfterExit=yes
