# build/master.project/ignition/units/swap.service
[Unit]
Description=Enable any swap files that are available
Requires=multi-user.target

RequiresMountsFor=/var

[Service]
Type=oneshot
RemainAfterExit=true

ExecStartPre=/bin/sh -c 'chmod 0600 -v ${Path}/*'
ExecStart=/bin/sh -c '/sbin/swapon -v ${Path}/*'
ExecStop=/bin/sh -c '/sbin/swapoff -v ${Path}/*'

[Install]
WantedBy=multi-user.target
