# build/master.project/salt/container/container-sync.service
[Unit]
Description=Container-Synchronization
Requires=container-load.service
After=container-load.service

[Service]
User=root
Type=oneshot
Environment="CONTAINER_DIR={{ container_path }}"
Environment="IMAGEDIR={{ container_image_path }}"
ExecStart=/bin/sh -c \
    "${CONTAINER_DIR}"/container-update.sh

[Install]
RequiredBy=container-sync.path
WantedBy=multi-user.target
