# build/master.project/salt/container/container-build.path
[Unit]
Description=Watch for any images that need to be built.

[Path]
PathChanged={{ container_path }}/build/
Unit=container-build.service

[Install]
WantedBy=multi-user.target
