# build/master.project/salt/container/container-build.service
[Unit]
Description=Container-Builder
Requires=network-online.target
After=network-online.target

[Service]
User=root
Type=oneshot
RemainAfterExit=yes
Environment="CONTAINER_DIR={{ container_path }}"
Environment="BUILDDIR={{ container_path }}/build"
Environment="IMAGEDIR={{ container_path }}/image"
Environment="TOOLDIR={{ container_path }}/tools"
ExecStart=/bin/sh -c \
    "${CONTAINER_DIR}"/container-build.sh

[Install]
RequiredBy=container-build.path
WantedBy=multi-user.target
