# build/master.project/salt/container/container-load.path
[Unit]
Description=Watch for any images that need to be loaded.

[Path]
PathChanged={{ container_path }}/image/
Unit=container-load.service

[Install]
WantedBy=multi-user.target
