# build/master.project/salt/container/container-load.service
[Unit]
Description=Container-Loader
Wants=container-build.service
After=container-build.service

[Service]
User=root
Type=oneshot
RemainAfterExit=yes
Environment="IMAGEDIR={{ container_image_path }}"
Environment="SERVICEDIR={{ container_service_path }}"
ExecStart=/bin/sh -c \
    "${SERVICEDIR}"/container-load.sh

[Install]
RequiredBy=container-load.path
WantedBy=multi-user.target
