# build/master.project/salt/container/container-build.service
[Unit]
Description=Container-Builder
Requires=network-online.target
After=network-online.target

[Service]
User=root
Type=oneshot
RemainAfterExit=yes
Environment="BUILDDIR={{ container_build_path }}"
Environment="IMAGEDIR={{ container_image_path }}"
Environment="SERVICEDIR={{ container_service_path }}"
ExecStart=/bin/sh -c \
    "${SERVICEDIR}"/container-build.sh

[Install]
RequiredBy=container-build.path
WantedBy=multi-user.target
