# build/master.project/salt/container/container-sync.path
[Unit]
Description=Watch for any images that need to be updated.

[Path]
PathChanged={{ container_image_path }}/
Unit=container-sync.service

[Install]
WantedBy=multi-user.target
