# build/master.project/salt/maintenance/command.service
[Unit]
Description={{ description }}

[Service]
Type=oneshot
ExecStart={{ command }}
