# build/master.project/salt/maintenance/command.timer
[Unit]
Description={{ description }}

# Please review the dropin for this timer to see the interval that
# this unit triggers its respective service.

[Install]
WantedBy=multi-user.target
